# localnotify/rootfs/etc/services.d/ localnotify/run
#!/usr/bin/execlineb -P
python3 /main.py